# systemd/netframe-predict.service
[Unit]
Description=NetFRAME predictive report (drive-failure + capacity forecast)
Documentation=file:/opt/netframe-monitor/netframe_predict.py
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
Nice=15
ExecStart=/usr/bin/python3 /opt/netframe-monitor/netframe_predict.py
TimeoutStartSec=900
