# systemd/netframe-backup.timer
[Unit]
Description=NetFRAME operational-memory backup timer (daily ~04:30)

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
